# deploy/systemd/image-to-text.service
[Unit]
Description=Image to Text — gunicorn web server
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/image-to-text
EnvironmentFile=/opt/image-to-text/.env
ExecStart=/opt/image-to-text/venv/bin/gunicorn \
    -w 2 \
    --worker-class gevent \
    -b 127.0.0.1:8000 \
    --timeout 300 \
    --keep-alive 5 \
    --log-level info \
    --access-logfile /opt/image-to-text/logs/access.log \
    --error-logfile /opt/image-to-text/logs/error.log \
    web.backend.app:app
Restart=always
RestartSec=5
TimeoutStartSec=180
StandardOutput=append:/opt/image-to-text/logs/gunicorn.log
StandardError=append:/opt/image-to-text/logs/gunicorn.log

[Install]
WantedBy=multi-user.target
